# _apache_spark/.hack/spark/appfs/etc/profile.d/_danny.sh
export HOME=/home/danny/.danny

export PATH=$HOME/bin:$PATH
